# 99-local.rules
ACTION=="add", KERNEL=="sdb*", SUBSYSTEM=="block", RUN+="/bin/sh -c 'echo /home/pi/Documents/4velma/import.sh | at now 2>/tmp/udev.txt'"
